# dotfiles/getty.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin viv --noclear %I $TERM
